# yt-streamer.service
[Unit]
Description=YouTube Live Streamer - Bulletin Rotation
After=network.target newsbot.service

[Service]
Type=simple
User=root
WorkingDirectory=/root/localaitv1
Environment="PATH=/root/localaitv1/env/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/localaitv1/env/bin/python yt_streamer.py

Restart=always
RestartSec=10
StartLimitInterval=0

StandardOutput=journal
StandardError=journal
SyslogIdentifier=yt-streamer

KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
